# config/systemd/user/hamr-mobile-restart.service
[Unit]
Description=Restart hamr-mobile after binary change

[Service]
Type=oneshot
ExecStart=systemctl --user restart hamr-mobile.service
