# config/systemd/user/hamr-mobile.service
[Unit]
Description=Hamr mobile web UI server
After=network-online.target

[Service]
Type=simple
# NOTE: the binary reads ~/.config/hamr + ~/.local/state/hamr/certs;
# dir migration is tracked as its own issue (code + data move together).
EnvironmentFile=%h/.config/hamr/env
ExecStart=%h/go/bin/hamr-mobile serve --address hamr.kingland.io:8444
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
