# config/systemd/user/hamr-ui-restart.path
[Unit]
Description=Watch hamr-ui binary for changes

[Path]
PathChanged=/home/drekar/go/bin/hamr-ui

[Install]
WantedBy=default.target
