# config/systemd/user/hamr-mobile-restart.path
[Unit]
Description=Watch hamr-mobile binary for changes

[Path]
PathChanged=/home/drekar/go/bin/hamr-mobile

[Install]
WantedBy=default.target
